# supervisord.conf
; Запускает MQTT-брокер и Airflow как два процесса внутри одного контейнера.

[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0

[program:mosquitto]
command=/usr/sbin/mosquitto -c /etc/mosquitto/mosquitto.conf
user=mosquitto
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; `airflow standalone` инициализирует БД и поднимает scheduler + webserver.
[program:airflow]
command=airflow standalone
user=airflow
autorestart=true
priority=20
environment=
    HOME="/home/airflow",
    AIRFLOW_HOME="/opt/airflow",
    PATH="/home/airflow/.local/bin:/usr/local/bin:/usr/bin:/bin",
    AIRFLOW__CORE__LOAD_EXAMPLES="false"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
